# deployment/congress-forms-dispatcher.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn

setuid ubuntu

chdir /vagrant/congress-forms-dispatcher

exec npm start
