# deployment/upstart-congress-forms-logger.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn

setuid root

exec remote_syslog -D --configfile=/vagrant/deployment/log_files.yml
