# deployment/upstart-congress-forms.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn

setuid ubuntu

chdir /vagrant

env PADRINO_ENV=production
exec /home/ubuntu/.rvm/wrappers/ruby-2.1.0@congress-forms/bundle exec rackup
