# packer_deployment/upstart-congress-forms.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn

setuid ubuntu

chdir /vagrant

exec /home/ubuntu/.rvm/wrappers/ruby-2.1.0@congress-forms/bundle exec rackup >> /vagrant/error.log 2>&1